# .config/systemd/user/swww-rand.timer
[Unit]
Description=Run swww-rand every 15 minutes

[Timer]
OnStartupSec=15m

OnUnitActiveSec=15m

[Install]
WantedBy=timers.target
