# .config/systemd/user/swww-rand.service
[Unit]
Description=swww randomizer

[Service]
Type=oneshot
ExecStart=%h/.local/bin/swww-rand
